# etc/upstart-startup.conf
# zfswatcher - ZFS monitoring daemon
#
# This file should be installed as /etc/init/zfswatcher.conf
# on architectures which use upstart(8) to bring up system
# daemons.

description	"ZFS monitoring daemon"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/sbin/zfswatcher
